# net-fs/samba/files/systemd/swat.socket
[Unit]
Description=Swat socket

[Socket]
ListenStream=901
Accept=yes

[Install]
WantedBy=sockets.target
