# net-fs/samba/files/systemd/swat@.service
[Unit]
Description=Swat daemon

[Service]
ExecStart=/usr/sbin/swat
StandardInput=socket
